# backend.ini
[program:fast-data-backend]
;user=nobody
command=/usr/bin/java -jar /fast.data.kcql.service-assembly-1.0.jar
redirect_stderr=true
stdout_logfile=/var/log/fast-data-backend.log
startretries=5
