# generator.ini
[program:generator]
;user=nobody
command=/usr/bin/java -jar /avro-kafka-generator-assembly-1.0.jar
redirect_stderr=true
stdout_logfile=/var/log/generator.log
startretries=5
